# live/management/ansible/roles/nextcloud_pg_dump/templates/nextcloud-pg-dump.service.j2
[Unit]
Description=Nextcloud Postgres logical dump (escape hatch — not the restore path)
Documentation=https://github.com/scottishwidow/scottishwidow-resources/blob/main/docs/runbooks/nextcloud-restore.md
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nextcloud-pg-dump.sh

SuccessExitStatus=0
