# live/management/ansible/roles/nextcloud_pg_dump/templates/nextcloud-pg-dump.timer.j2
[Unit]
Description=Daily Nextcloud Postgres logical dump

[Timer]
OnCalendar={{ nextcloud_pg_dump_oncalendar }}

Persistent=true

RandomizedDelaySec=120

[Install]
WantedBy=timers.target
